# osx/crontab
env EDITOR=nano crontab -e

*/10 * * * * /Users/sorokan6/sync.sh
